# supervisord.conf
[supervisord]
nodaemon=true

[program:asterisk]
command=/usr/sbin/asterisk -f -U asterisk -G asterisk -vvvg -c
redirect_stderr=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
priority=1